# backup.timer
[Unit]
Description=Run Automatic LXD Backup
 
[Timer]
OnCalendar=*-*-* ${TIME}:00:00
Persistent=true
 
[Install]
WantedBy=multi-user.target
